# algo-trade-platform/server/algo-trading-pipeline.service
[Unit]
Description=Algo Trading Pipeline Service
After=network.target postgresql.service algo-trading-strategy.service
Requires=algo-trading-strategy.service

[Service]
User=root
WorkingDirectory=/path/to/algo-trade-platform
ExecStart=/usr/bin/python3 /path/to/algo-trade-platform/backend/app/pipeline/pipeline_manager.py
Restart=always
RestartSec=10
Environment=PIPELINE_CONFIG_PATH=/path/to/algo-trade-platform/config/pipeline_config.json
Environment=DB_PASSWORD=password
Environment=PYTHONPATH=/path/to/algo-trade-platform

[Install]
WantedBy=multi-user.target
