# algo-trade-platform/server/algo-trading-strategy.service
[Unit]
Description=Algo Trading Strategy Service
After=network.target postgresql.service

[Service]
User=root
WorkingDirectory=/path/to/algo-trade-platform
ExecStart=/usr/bin/python3 /path/to/algo-trade-platform/backend/app/strategies/Stocks/attempt.py
Restart=always
RestartSec=10
Environment=DB_PASSWORD=password
Environment=CONFIG_PATH=/path/to/algo-trade-platform/config/config.json
Environment=PYTHONPATH=/path/to/algo-trade-platform

[Install]
WantedBy=multi-user.target
