# algo-trade-platform/server/ibkr-streamdata.service
[Unit]
Description=IBKR Streamdata Service with Telegram Notifications
After=network.target
Wants=network.target

[Service]
Type=simple
User=anurag
WorkingDirectory=/home/anurag/Desktop/Algo-Trading/algo-trade-platform
ExecStart=/home/anurag/Desktop/Algo-Trading/algo-trade-platform/backend/app/data/venv/bin/python /home/anurag/Desktop/Algo-Trading/algo-trade-platform/backend/app/data/start_streamdata.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
Environment=TELEGRAM_BOT_TOKEN=8161925605:AAEFvQtiTmnjUnLiR1BAfVnxVfu8aOg-Yhk
Environment=TELEGRAM_CHAT_ID=US-stocks.Jr
Environment=PYTHONPATH=/home/anurag/Desktop/Algo-Trading/algo-trade-platform
Environment=VIRTUAL_ENV=/home/anurag/Desktop/Algo-Trading/algo-trade-platform/backend/app/data/venv
Environment=PATH=/home/anurag/Desktop/Algo-Trading/algo-trade-platform/backend/app/data/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target 